# jw-dioxus.service
[Unit]
Description=JW Dioxus Web Application
After=network.target
Wants=network.target

[Service]
Type=simple
User=lc
Group=lc
WorkingDirectory=/home/lc/jw
ExecStart=dx serve --release
Restart=always
RestartSec=1
StandardOutput=journal
StandardError=journal
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=RUST_LOG=info

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/lc/jw
PrivateTmp=true

[Install]
WantedBy=multi-user.target
